# mtk-sdk-wifi/files/etc/uci-defaults/50-mt7688-iwinfo
#!/bin/sh

[ -f /usr/bin/iwinfo.rt ] && mv /usr/bin/iwinfo.rt /usr/bin/iwinfo
[ -f /lib/libiwinfo.so ] && mv /usr/lib/libiwinfo.so

return 0
